# noiseforwarder/src/dist/noiseforwarder.service
# copy this to /etc/systemd/system
# then sudo systemctl daemon-reload
# and enable with sudo systemctl enable --now noiseforwarder

[Unit]
Description=noiseforwarder
After=network.target

[Service]
WorkingDirectory=/home/bertrik/noiseforwarder
ExecStart=/home/bertrik/noiseforwarder/bin/noiseforwarder
Type=simple
KillMode=process
Restart=on-failure
RestartSec=10s
TimeoutStopSec=15s
User=bertrik
Environment="NOISEFORWARDER_OPTS=-Xmx64M"

[Install]
WantedBy=multi-user.target
